# deploy/projects-hub.service
[Unit]
Description=Projects Hub (Next.js)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/projects-hub
ExecStart=/usr/bin/env pnpm start
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=HOSTNAME=127.0.0.1
EnvironmentFile=-/opt/projects-hub/.env

[Install]
WantedBy=multi-user.target
